# idrive/supervisord.conf
[supervisord]
nodaemon=true

[program:idrive]
command=/usr/bin/idrive_admin.py start
